# system/hysteria2-manager.service
# ==================== hysteria2-manager.service ====================
# Hysteria2 Manager WebUI Service
# Location: /etc/systemd/system/hysteria2-manager.service
# Description: Manages the web interface and API for Hysteria2 proxy

[Unit]
Description=Hysteria2 Manager - Web Management Interface
Documentation=https://github.com/1439616687/hysteria2-manager
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

# Dependencies
Requires=network.target

[Service]
# Service Type
Type=simple

# User and Group
User=root
Group=root

# Working Directory
WorkingDirectory=/opt/hysteria2-manager

# Environment Variables
Environment="PATH=/opt/hysteria2-manager/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/hysteria2-manager"
Environment="PYTHONUNBUFFERED=1"
Environment="LANG=en_US.UTF-8"
Environment="LC_ALL=en_US.UTF-8"

# Pre-start checks
ExecStartPre=/bin/bash -c 'if [ ! -d /opt/hysteria2-manager ]; then echo "Installation directory not found"; exit 1; fi'
ExecStartPre=/bin/bash -c 'if [ ! -f /opt/hysteria2-manager/data/config.json ]; then echo "Configuration file not found"; exit 1; fi'
ExecStartPre=/bin/bash -c 'if [ ! -f /opt/hysteria2-manager/venv/bin/python ]; then echo "Python virtual environment not found"; exit 1; fi'

# Start Command
ExecStart=/opt/hysteria2-manager/venv/bin/python /opt/hysteria2-manager/hysteria2_manager.py

# Restart Policy
Restart=always
RestartSec=10
StartLimitInterval=120
StartLimitBurst=5

# Process Management
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
TimeoutStartSec=30
TimeoutStopSec=10

# Resource Limits
LimitNOFILE=65535
LimitNPROC=4096
LimitCORE=0
# Memory limit (optional, adjust as needed)
# MemoryLimit=256M
# MemoryAccounting=true
# CPUQuota=50%

# Security Settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hysteria2-manager/data /var/log/hysteria2
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true
# RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# RestrictNamespaces=true
# RestrictRealtime=true
# SystemCallFilter=@system-service

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hysteria2-manager

# Performance Tuning
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=4

# Health Check (optional)
# Watchdog=30s
# WatchdogSec=30s

[Install]
WantedBy=multi-user.target

# ==================== hysteria2-client.service ====================
# Hysteria2 Client Service
# Location: /etc/systemd/system/hysteria2-client.service
# Description: Hysteria2 proxy client that creates TUN interface

[Unit]
Description=Hysteria2 Client - Proxy Connection Service
Documentation=https://v2.hysteria.network/
Documentation=https://github.com/apernet/hysteria
After=network-online.target nss-lookup.target
Wants=network-online.target

# Dependency on manager service
BindsTo=hysteria2-manager.service
After=hysteria2-manager.service

# Network requirements
Requires=network.target
After=network.target

[Service]
# Service Type
Type=simple

# User and Group (requires root for TUN device)
User=root
Group=root

# Working Directory
WorkingDirectory=/etc/hysteria2

# Environment Variables
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LANG=en_US.UTF-8"
Environment="LC_ALL=en_US.UTF-8"

# Pre-start Operations
# Clean up any existing TUN interface
ExecStartPre=-/usr/sbin/ip link delete hytun 2>/dev/null
# Wait for network to be fully ready
ExecStartPre=/bin/sleep 2
# Verify configuration file exists
ExecStartPre=/bin/bash -c 'if [ ! -f /etc/hysteria2/client.yaml ]; then echo "Client configuration not found. Please select a node first."; exit 1; fi'
# Test configuration (optional, comment out if causes issues)
# ExecStartPre=/usr/local/bin/hysteria client -c /etc/hysteria2/client.yaml --test

# Main Process
ExecStart=/usr/local/bin/hysteria client -c /etc/hysteria2/client.yaml

# Post-stop Operations
# Clean up TUN interface
ExecStopPost=-/usr/sbin/ip link delete hytun 2>/dev/null
# Flush route cache
ExecStopPost=-/bin/bash -c 'ip route flush cache'
# Optional: Clear iptables rules
# ExecStopPost=-/bin/bash -c 'iptables -t nat -F'

# Restart Policy
Restart=on-failure
RestartSec=10
StartLimitInterval=120
StartLimitBurst=5

# Process Management
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
TimeoutStartSec=30
TimeoutStopSec=10

# Resource Limits
LimitNOFILE=1048576
LimitNPROC=512
LimitCORE=0
# Memory limit (optional)
# MemoryLimit=512M
# MemoryAccounting=true
# CPUQuota=80%

# Network Capabilities (IMPORTANT!)
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

# Security Settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/hysteria2 /var/log/hysteria2
# Allow access to /dev/net/tun device
PrivateDevices=no
DeviceAllow=/dev/net/tun rw
# ProtectKernelTunables=false  # Need to modify network settings
# ProtectKernelModules=false   # Need to load TUN module
# RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET
# SystemCallFilter=@system-service @network-io @privileged

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hysteria2-client

# Performance Tuning
Nice=-10
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
IOSchedulingClass=best-effort
IOSchedulingPriority=2

# Network Settings
# Allow binding to privileged ports
# PrivateNetwork=no
# IPAddressDeny=
# IPAddressAllow=any

# Environment Checks (optional)
# AssertPathExists=/usr/local/bin/hysteria
# AssertPathExists=/etc/hysteria2
# AssertFileNotEmpty=/etc/hysteria2/client.yaml
# AssertCapability=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
# Also start when manager starts
# WantedBy=hysteria2-manager.service
