# system/hysteria2-client.service
[Unit]
Description=Hysteria2 Client Service
Documentation=https://v2.hysteria.network/
After=network-online.target nss-lookup.target
Wants=network-online.target

# 依赖检查
ConditionPathExists=/usr/local/bin/hysteria
ConditionPathExists=/etc/hysteria2/client.yaml

# 冲突检查
Conflicts=hysteria-server.service
Conflicts=v2ray.service
Conflicts=xray.service

[Service]
# 服务类型
Type=simple
Restart=on-failure
RestartSec=10

# 用户和权限
User=root
Group=root

# 工作目录
WorkingDirectory=/etc/hysteria2

# 环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HYSTERIA_LOG_LEVEL=info"
Environment="GOGC=50"

# 网络命名空间（可选）
# NetworkNamespacePath=/var/run/netns/hysteria

# TUN接口清理和准备
ExecStartPre=-/usr/sbin/ip link delete hytun 2>/dev/null
ExecStartPre=/bin/sleep 2
ExecStartPre=/bin/bash -c 'test -f /etc/hysteria2/client.yaml || exit 1'
ExecStartPre=/bin/bash -c '/usr/local/bin/hysteria version || exit 1'

# 主进程启动命令
ExecStart=/usr/local/bin/hysteria client -c /etc/hysteria2/client.yaml

# 停止后清理
ExecStopPost=-/usr/sbin/ip link delete hytun 2>/dev/null
ExecStopPost=/bin/bash -c 'iptables -t nat -D POSTROUTING -o hytun -j MASQUERADE 2>/dev/null || true'
ExecStopPost=/bin/bash -c 'ip rule del fwmark 0x1 table 100 2>/dev/null || true'
ExecStopPost=/bin/bash -c 'ip route del default dev hytun table 100 2>/dev/null || true'

# 重载配置
ExecReload=/bin/kill -HUP $MAINPID

# 进程管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=30
TimeoutStopSec=10
SendSIGKILL=yes

# 重启策略
StartLimitInterval=120
StartLimitBurst=5

# 资源限制
LimitNOFILE=1048576
LimitNPROC=512
LimitCORE=infinity
TasksMax=512

# 内存管理
MemoryMax=256M
MemorySwapMax=0
MemoryAccounting=true

# CPU管理
CPUWeight=200
CPUAccounting=true

# 网络权限
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SETUID CAP_SETGID
SecureBits=keep-caps

# 安全设置
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=false
LockPersonality=true
RemoveIPC=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET
SystemCallFilter=~@clock @cpu-emulation @debug @module @mount @obsolete @raw-io @reboot @swap

# 文件系统权限
ReadWritePaths=/etc/hysteria2
ReadWritePaths=/var/log/hysteria2
ReadWritePaths=/run
ReadWritePaths=/var/run
ReadOnlyPaths=/usr/local/bin

# 设备访问（TUN需要）
PrivateDevices=no
DeviceAllow=/dev/net/tun rw
DeviceAllow=/dev/null rw
DeviceAllow=/dev/urandom r
DevicePolicy=strict

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hysteria2-client
LogLevelMax=info

# 环境文件（可选）
EnvironmentFile=-/etc/default/hysteria2-client

# 监控
WatchdogSec=60
NotifyAccess=none

# OOM设置
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target
Alias=hy2-client.service
Also=hysteria2-manager.service
